# bin/robogames.conf
description     "Robogames!"

start on startup

chdir /root/robogames
respawn
respawn limit 5 5 
exec node app.js >> log.txt 2>&1